# scripts/crontab
*/2 * * * * /opt/nodeprojects/GrosvenorLocation/scripts/checkRestartTwitterStream.sh >> /var/log/GrosvenorLocation/shellRestartTwitterStream.log
*/1 * * * * /opt/nodeprojects/GrosvenorLocation/scripts/moveFiles.sh >> /var/log/GrosvenorLocation/shellmoveFiles.log
30 23 * * * /home/GvAdmin/Hive_Partition.sh >> /var/log/hivepartitionscreation.log
